# deploy/systemd/web-admin.service
[Unit]
Description=Web Admin Panel (Production)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/web
ExecStart=/opt/web/venv/bin/python main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
